# deploy/com.polymarket.control-plane.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.polymarket.control-plane</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/joez/Desktop/polymarket-mvp/.venv311/bin/python</string>
    <string>-m</string>
    <string>polymarket_mvp.control_plane</string>
    <string>--port</string>
    <string>8788</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/joez/Desktop/polymarket-mvp</string>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/joez/Desktop/polymarket-mvp/var/log/control-plane.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/joez/Desktop/polymarket-mvp/var/log/control-plane.err</string>
</dict>
</plist>
