# deploy/com.polymarket.autopilot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.polymarket.autopilot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/joez/Desktop/polymarket-mvp/.venv311/bin/python</string>
    <string>-m</string>
    <string>polymarket_mvp.autopilot</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/joez/Desktop/polymarket-mvp</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/joez/Desktop/polymarket-mvp/var/log/autopilot.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/joez/Desktop/polymarket-mvp/var/log/autopilot.err</string>
</dict>
</plist>
